# deploy/systemd/momo.service
[Unit]
Description=MoMo Wireless Security Platform
Documentation=https://github.com/M0M0Sec/MoMo
After=network-online.target momo-ap.service
Wants=network-online.target

# Don't start until first boot wizard completes
ConditionPathExists=|/etc/momo/.setup_complete
ConditionPathExists=|/etc/momo/momo.yml

[Service]
Type=simple
WorkingDirectory=/opt/momo

# Use /etc/momo/momo.yml as primary config location
ExecStart=/opt/momo/.venv/bin/momo run -c /etc/momo/momo.yml

# Environment files
EnvironmentFile=-/etc/default/momo
EnvironmentFile=-/etc/systemd/system/momo.service.d/env.conf

# Restart policy
Restart=always
RestartSec=5
LimitNOFILE=65535

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=momo

[Install]
WantedBy=multi-user.target


