# deploy/firstboot/momo-firstboot.service
# ==============================================================================
# MoMo First Boot Wizard Service
# ==============================================================================
# This systemd service handles the first boot setup wizard.
#
# Installation:
#   sudo cp momo-firstboot.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable momo-firstboot.service
#
# The service runs once on boot and:
#   - If /etc/momo/.setup_complete exists: Exits immediately (normal boot)
#   - If /boot/momo-config.yml exists: Applies headless config
#   - Otherwise: Starts WiFi AP and web wizard
# ==============================================================================

[Unit]
Description=MoMo First Boot Setup Wizard
Documentation=https://github.com/M0M0Sec/MoMo

# Run before main MoMo service
Before=momo.service
Before=momo-web.service

# Wait for system to be fully ready
After=network.target
After=systemd-resolved.service
After=NetworkManager.service
After=multi-user.target

# Don't start if setup is complete (check multiple locations)
ConditionPathExists=!/etc/momo/.setup_complete
ConditionPathExists=!/opt/momo/configs/.setup_complete

[Service]
Type=simple

# Run as root (needed for network configuration)
User=root
Group=root

# Working directory
WorkingDirectory=/opt/momo

# Python path - use virtual environment
Environment=PYTHONPATH=/opt/momo
Environment=LANG=en_US.UTF-8
Environment=PATH=/opt/momo/.venv/bin:/usr/local/bin:/usr/bin:/bin

# Pre-start: Wait for system stability and release wlan0
ExecStartPre=/bin/sleep 5
ExecStartPre=-/usr/bin/nmcli dev set wlan0 managed no
ExecStartPre=-/usr/bin/nmcli dev disconnect wlan0
ExecStartPre=-/usr/bin/killall wpa_supplicant

# Main command - use venv python
ExecStart=/opt/momo/.venv/bin/python -m momo.firstboot.entry

# Timeout (allow plenty of time for wizard)
TimeoutStartSec=3600

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=momo-firstboot

[Install]
WantedBy=multi-user.target

